# deploy/cron/cleanup-old-mail.cron
0 * * * * root /opt/temp-mail/scripts/cleanup_old_mail.sh >/dev/null 2>&1
